# buildroot/raspberrypi0/ovl/etc/profile.d/custom.sh
# Acces aux scripts personnels et code metier.
PATH=$PATH:/usr/local/bin

# Prompt indiquant nom d'hote et repertoire courant.
PS1='\h[\W]\$ '
